# weatherhat.service
[Unit]
Description=Weather HAT MQTT Publisher (Improved)
Documentation=https://github.com/pimoroni/weather-station
After=network-online.target time-sync.target
Wants=network-online.target time-sync.target

# Optional: if using lgpiod daemon
# After=lgpiod.service

[Service]
Type=simple
User=weather
Group=weather

# Working directory
WorkingDirectory=/home/weather/weather-station/bin

# Load environment configuration from file
# Create this file from mqtt.env.example
EnvironmentFile=/home/weather/weather-station/config/mqtt.env

# Environment setup - use virtual environment's Python directly
Environment="PATH=/home/weather/.virtualenvs/pimoroni/bin:/usr/local/bin:/usr/bin:/bin"
Environment="VIRTUAL_ENV=/home/weather/.virtualenvs/pimoroni"

# Give hardware time to initialize after boot
ExecStartPre=/bin/sleep 15

# Run the improved MQTT publisher script with structured logging
ExecStart=/home/weather/.virtualenvs/pimoroni/bin/python /home/weather/weather-station/bin/mqtt-publisher.py

# Restart configuration - robust for remote deployment
# Generous limits: transient I2C/hardware glitches can cause several
# rapid restarts. 10 restarts per hour prevents permanent death while
# still catching genuine boot loops.
Restart=on-failure
RestartSec=60
TimeoutStopSec=15
StartLimitBurst=10
StartLimitIntervalSec=3600

# Resource limits - improved version uses less CPU
MemoryMax=128M
CPUQuota=25%

# Logging - structured logs work well with journald
StandardOutput=journal
StandardError=journal
SyslogIdentifier=weatherhat

# Security hardening (optional)
# ProtectSystem=strict
# ProtectHome=read-only
# NoNewPrivileges=yes
# PrivateTmp=yes

[Install]
WantedBy=multi-user.target
