# weatherhat.service
[Unit]
Description=Weather HAT MQTT Publisher
Documentation=https://github.com/timFinn/weatherhat-python

# Wait for network and time sync - critical for remote stations
After=network-online.target time-sync.target
Wants=network-online.target time-sync.target

# Optional: if using lgpiod daemon
# After=lgpiod.service

[Service]
Type=simple
User=garden
Group=garden

# Working directory
WorkingDirectory=/home/garden/weatherhat-python/examples

# Environment setup - use virtual environment's Python directly
Environment="PATH=/home/garden/.virtualenvs/pimoroni/bin:/usr/local/bin:/usr/bin:/bin"
Environment="VIRTUAL_ENV=/home/garden/.virtualenvs/pimoroni"

# Give hardware time to initialize after boot
ExecStartPre=/bin/sleep 15

# Run the MQTT publisher script
ExecStart=/home/garden/.virtualenvs/pimoroni/bin/python /home/garden/weatherhat-python/examples/mqtt.py

# Restart configuration - robust for remote deployment
Restart=on-failure
RestartSec=30
StartLimitBurst=5
StartLimitIntervalSec=300

# Watchdog - restart if script hangs (optional, requires code support)
# WatchdogSec=120

# Resource limits appropriate for Pi Zero
MemoryMax=128M
CPUQuota=50%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=weatherhat

[Install]
WantedBy=multi-user.target
