# weatherhat-display.service
[Unit]
Description=Weather HAT Display (Power Saving)
Documentation=https://github.com/pimoroni/weather-station
After=weatherhat.service
# Optional: start after main weather service

# Restart rate-limit: stop retrying after 3 failures in 5 minutes.
StartLimitBurst=3
StartLimitIntervalSec=300

[Service]
Type=simple
User=weather
Group=weather

# Working directory
WorkingDirectory=/home/weather/weather-station/bin

# Environment setup
Environment="PATH=/home/weather/.virtualenvs/pimoroni/bin:/usr/local/bin:/usr/bin:/bin"
Environment="VIRTUAL_ENV=/home/weather/.virtualenvs/pimoroni"
Environment="TEMP_OFFSET=-7.5"
Environment="DISPLAY_SLEEP_TIMEOUT=30"

# Give hardware time to initialize
ExecStartPre=/bin/sleep 5

# Run the display script
ExecStart=/home/weather/.virtualenvs/pimoroni/bin/python /home/weather/weather-station/bin/display-interface.py

# Restart on failure
Restart=on-failure
RestartSec=10

# Resource limits - display uses minimal resources when sleeping
MemoryMax=128M
CPUQuota=15%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=weatherhat-display

[Install]
WantedBy=multi-user.target
